# supervisor-conf.d/run_server_power.conf
[program:run_server_power]
command = /home/seina/homedisplay-env/bin/python run.py
directory = /home/seina/home-info-display/server_power_publisher
user=seina

