# supervisor-conf.d/run_indoor_quality.conf
[program:indoor_quality]
command = /home/seina/homedisplay-env/bin/python run.py
directory = /home/seina/home-info-display/indoor_quality_publisher
user=seina

