# supervisor-conf.d/run_ping.conf
[program:run_ping]
command = /home/seina/homedisplay-env/bin/python run.py
directory = /home/seina/home-info-display/ping_publisher
user=seina

