# supervisor-conf.d/run_display_control.conf
[program:run_display_control]
command = /home/seina/homedisplay-env/bin/python run.py
directory = /home/seina/home-info-display/display_control_consumer
user=seina

